# roles/vault/templates/etc_systemd_system_tailscale-cert-vault.timer.j2
[Unit]
Description="Generate a tailscale cert for hashicorp vault"

[Timer]
OnBootSec=120s
OnUnitActiveSec=1w

[Install]
WantedBy=multi-user.target
