# roles/vault/templates/etc/systemd/system/vault.service.d/override.conf
# {{ ansible_managed }}

[Service]
# lets systemd create directory `/run/vault`
RuntimeDirectory=vault
